# systemd/docker-media.service
[Unit]
Description=Docker Media Stack (gluetun + qbittorrent)
Wants=network-online.target
After=network-online.target docker.service snap.docker.dockerd.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/docker-media
ExecStart=/usr/bin/docker compose up -d --remove-orphans
ExecStop=/usr/bin/docker compose down
Restart=on-failure
RestartSec=15
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target
